# cron.d/pull_pause_container
0,15,30,45 *    * * * root /usr/local/share/eks/pull-sandbox-image.sh > /tmp/pause-container-pull.log
